# bcm-monitoring/bcm-metrics-server.service
[Unit]
Description=BCM Process Metrics HTTP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /tmp/metrics-server.py 9256
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
